# deploy/systemd/remotellm-connector.service
[Unit]
Description=RemoteLLM Connector - Bridge local LLM to external API
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=remotellm
Group=remotellm
WorkingDirectory=/opt/remotellm
ExecStart=/opt/remotellm/venv/bin/python -m remotellm.connector
Restart=always
RestartSec=10

# Environment file for configuration
EnvironmentFile=/etc/remotellm/connector.env

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadOnlyPaths=/

# Allow writing to logs if needed
ReadWritePaths=/var/log/remotellm

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=remotellm-connector

[Install]
WantedBy=multi-user.target
