# deploy/launchd/com.remotellm.broker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.remotellm.broker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/remotellm-broker</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>REMOTELLM_BROKER_HOST</key>
        <string>0.0.0.0</string>
        <key>REMOTELLM_BROKER_PORT</key>
        <string>8443</string>
        <!-- Comma-separated list of valid connector tokens -->
        <key>REMOTELLM_BROKER_CONNECTOR_TOKENS</key>
        <string>token1,token2</string>
        <!-- Comma-separated list of valid user API keys -->
        <key>REMOTELLM_BROKER_USER_API_KEYS</key>
        <string>user-key-1,user-key-2</string>
        <!-- Path to YAML config with connector token → llm_api_key mappings -->
        <key>REMOTELLM_BROKER_CONNECTOR_CONFIG</key>
        <string>/usr/local/etc/remotellm/connectors.yaml</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/usr/local/var/log/remotellm/broker.log</string>

    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/remotellm/broker.error.log</string>

    <key>WorkingDirectory</key>
    <string>/usr/local/opt/remotellm</string>
</dict>
</plist>
